# com.bi-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd 配置：保活 bi-bridge，崩溃/退出后自动重启。
  安装（把下面占位 __REPO__ / __ANALYST__ 改成你的真实值后）：
    cp com.bi-bridge.plist ~/Library/LaunchAgents/com.bi-bridge.me.plist
    # 编辑里面的 __REPO__（仓库根的绝对路径）、Label 改成带你名字的唯一值（多机器人时避免冲突）
    launchctl load   ~/Library/LaunchAgents/com.bi-bridge.me.plist   # 启动
    launchctl unload ~/Library/LaunchAgents/com.bi-bridge.me.plist   # 优雅停（发 SIGTERM，不会 kill -9，订阅不泄漏）

  注意：
  - launchd 默认 StandardInput=null，会让 consume 立即收到 EOF。本脚本用 --timeout
    循环重启、不依赖常驻 stdin，所以 EOF 无害（每轮 consume 由 --timeout 控制时长）。
  - KeepAlive=true 保证脚本意外退出后被拉起。
  - 务必让 PATH 能找到 lark-cli 和 claude（见 EnvironmentVariables.PATH）。
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bi-bridge.__ANALYST__</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REPO__/bi-bridge.sh</string>
        <string>__REPO__/config.__ANALYST__.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <!-- 改成你机器上能找到 lark-cli / claude / jq 的实际 PATH -->
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__REPO__/.bi-bridge/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO__/.bi-bridge/launchd.err.log</string>

    <!-- 崩溃后重启的最小间隔，防抖 -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
